# centos_nginx_php5/nginx_php.ini
[program:nginx]
command=/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
user=root

[program:php]
command=/sbin/php-fpm
autostart=true
autorestart=true
user=root

